# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:backend]
command=/bin/sh /app/start-backend.sh
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="5000",DATABASE_URL="%(ENV_DATABASE_URL)s",IMAP_HOST="%(ENV_IMAP_HOST)s",IMAP_USER="%(ENV_IMAP_USER)s",IMAP_PASSWORD="%(ENV_IMAP_PASSWORD)s",IMAP_ALLOWED_SENDER="%(ENV_IMAP_ALLOWED_SENDER)s",EMAIL_POLL_INTERVAL="%(ENV_EMAIL_POLL_INTERVAL)s"
priority=20
